# ops/com.claw.chrome-cdp.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.claw.chrome-cdp</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/hector/Projects/Dr.-strange/ops/chrome-cdp-launcher.sh</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>CLAW_CHROME_PORT</key>
    <string>9250</string>
    <key>CLAW_CHROME_PROFILE_DIR</key>
    <string>/Users/hector/.claw/chrome-profile</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/hector/.claw/chrome-cdp.stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/hector/.claw/chrome-cdp.stderr.log</string>
</dict>
</plist>
